# supervisored/supervisord.ini
[supervisord]
nodaemon=true

[program:postgres]
command=/usr/bin/postgres_service.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/postgres_err.log
stdout_logfile=/var/log/postgres_out.log
stopsignal=QUIT


[program:tomcat]
command=/opt/tomcat/latest/bin/catalina.sh run
autorestart=true
startsecs=20
stopsignal=INT
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/catalina.out
stderr_logfile=/var/log/catalina.out
environment=JAVA_HOME="/usr/lib/jvm/default-jvm",JAVA_BIN="/usr/lib/jvm/default-jvm/bin"